# ghostlink-controller.service.j2
[Unit]
Description=GhostLink Controller (read-only, metrics)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ghostlink
Group=ghostlink
EnvironmentFile=/etc/default/ghostlink
WorkingDirectory=/opt/ghostlink
ExecStart=/usr/bin/python3 /opt/ghostlink/gl_controller_metrics.py
Restart=on-failure
NoNewPrivileges=yes
ProtectSystem=full
ProtectHome=true
PrivateTmp=true
ProtectKernelTunables=yes
ProtectControlGroups=yes
LockPersonality=yes
MemoryDenyWriteExecute=true

[Install]
WantedBy=multi-user.target
